# shell/crontab
# 每周一到周五的下午5:00发送夕阳红
0 17 * * 1-5 /home/chstockbot/shell/sendxyh.sh